# 02-systemd_desktop/etc,systemd/system/morty.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/morty -listen 127.0.0.1:3000 -key d2hhdGV2ZXI= -ipv6 -socks5 127.0.0.1:9050
