# 02-systemd_laptop/etc,systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin jeli --noclear %I $TERM
