# 02-systemd/HOME,.config,systemd/user/mpdas.service
[Unit]
Description=AudoScrobbler client for MPD

[Service]
ExecStart=/usr/bin/mpdas -c %h/.config/mpdasrc

[Install]
WantedBy=default.target
