# 02-systemd/HOME,.config,systemd/user/searx.service
[Unit]
Description=Basic local Searx instance
After=tor.service

[Service]
ExecStart=/usr/local/bin/searx-run

[Install]
WantedBy=default.target
